# api/ansible/roles/redis/templates/laravel-worker.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:listen redis --tries=5 --timeout=3600 --queue=import,default,year_end_adjustment,employee,delivery
directory=/vagrant
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/vagrant/storage/logs/worker.log
